# usb_autocopy.service
[Install]
WantedBy=multi-user.target

[Service]
WorkingDirectory=/mnt/reserved/kurasushi_handwash
ExecStart=python3 usb_autocopy.py
Restart=always
RestartSec=1
